# packaging/indicator.path
[Unit]
Description=This path will tell you that indicator is ready to launch

[Path]
PathExists=/run/.wm_ready

